# references/meshtastic-detection.service
[Unit]
Description=Meshtastic Detection Sensor USB Receiver
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/venv/bin/python scripts/usb_receiver.py --port __SERIAL_PORT__
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
